# artemis.service
[Unit]
Description=Artemis External Pentest Automation Suite (Web)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/artemis
ExecStart=/usr/bin/python3 /opt/artemis/artemis_web.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# ── Credentials — edit these before starting ──────────────────────────
# Or set them in /etc/artemis.env and uncomment EnvironmentFile below
Environment="ARTEMIS_USER=artemis"
Environment="ARTEMIS_PASSWORD=changeme123!"
Environment="ARTEMIS_SECRET=replace-with-a-long-random-string"
# EnvironmentFile=/etc/artemis.env

# ── Runtime environment ───────────────────────────────────────────────
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/go/bin"
Environment="GOPATH=/root/go"
Environment="HOME=/root"

[Install]
WantedBy=multi-user.target
